# ansible/configs/supervisord.ini
[program:gunicorn]
command=/opt/dss/venv/bin/gunicorn dss:app -c /vagrant/gunicorn.conf.py
directory=/opt/dss/serverRoot
user=vagrant
autostart=true
autorestart=true
redirect_stderr=true
